# webscraper.service
[Unit]
Description=Web Scraper Selenium Flask App
After=network.target xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/webscraper
Environment=DISPLAY=:99
Environment=HOME=/home/ubuntu
ExecStart=/opt/webscraper/start_app.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target